# tomcat7/tomcat.conf
[program:tomcat]
command=/opt/tomcat/bin/catalina.sh run
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/tomcat/logs/catalina.out
