# python/scribed.conf
[program:error_to_scribed]
command=sh -c "tail -F /opt/logs/error.log | /opt/scribed/tools/scribe_tail -s $SCRIBE_IP:1463 $APP_NAME_error"
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/logs/error_to_scribed.log

[program:std_to_scribed]
command=sh -c "tail -F /opt/logs/std.log | /opt/scribed/tools/scribe_tail -s $SCRIBE_IP:1463 $APP_NAME_std"
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/logs/std_to_scribed.log

[group:scribed]
programs = error_to_scribed, std_to_scribed
