# nodejs/nodejs.conf
[program:nodejs]
command=node --harmony app.js
directory=/opt/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/logs/std.log
