# php/nginxphp.conf
[program:nginxphp]
command=sh -c "php-fpm && nginx"
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/logs/nginxphp.log
